# ops/scheduler/systemd/breathwise-pipeline-standard.timer
[Unit]
Description=Breathwise TR standard pipeline timer

[Timer]
OnCalendar=*-*-* 02:15:00
Persistent=true
Unit=breathwise-pipeline@standard.service

[Install]
WantedBy=timers.target

